# flask-app.service
[Unit]
Description=Finnegans Wake Flask App

[Service]
ExecStart=/app/startup
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
