# flask-app.service
[Unit]
Description=Systemd Config For Flask App Launch & Restart

[Service]
ExecStart=/app/startup
